# provisioning/roles/ehacking-deploy/templates/ehacking.service.j2
[Unit]
Description=eHacking offline training stack (rootless podman compose)
# The stack talks to the rootless podman API socket; make sure it's up first.
After=podman.socket
Wants=podman.socket

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory={{ ehacking_dir }}
ExecStartPre={{ ehacking_dir }}/bin/make-flags.sh
ExecStart={{ ehacking_dir }}/bin/compose up -d
ExecStop={{ ehacking_dir }}/bin/compose down

[Install]
WantedBy=default.target
